# ubuntu/config/rclone-gdrive.service
[Unit]
Description=Rclone Mount - Google Drive (udugunmur)
Documentation=https://rclone.org/commands/rclone_mount/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=edugonmor
Group=edugonmor

# Configuración de rclone
Environment=RCLONE_CONFIG=/home/edugonmor/.config/rclone/rclone.conf

# Comando de montaje
ExecStart=/usr/bin/rclone mount gdrive-udugunmur: /mnt/disk2/rclone/gdrive \
    --config=${RCLONE_CONFIG} \
    --vfs-cache-mode writes \
    --vfs-cache-max-age 24h \
    --vfs-read-chunk-size 64M \
    --vfs-read-chunk-size-limit 1G \
    --buffer-size 64M \
    --dir-cache-time 72h \
    --poll-interval 15s \
    --log-level INFO \
    --log-file /var/log/rclone-gdrive.log \
    --allow-other \
    --async-read

# Desmontaje limpio
ExecStop=/bin/fusermount -uz /mnt/disk2/rclone/gdrive

# Reintentos
Restart=on-failure
RestartSec=10

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
